# upstart/spaceifyappman.conf
#!upstart
description "Spaceify Application Manager Server"
author      "Spaceify Oy"

console none

start on (local-filesystems and started spaceify)
stop on shutdown

respawn

script
	export HOME="/root"

	echo $$ > /var/run/spaceifyappman.pid

	exec /var/lib/spaceify/data/scripts/start_appman.sh
end script

pre-start script
	echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/lib/spaceify/data/logs/spaceifyappman.log
end script

pre-stop script
	rm /var/run/spaceifyappman.pid
	echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/lib/spaceify/data/logs/spaceifyappman.log
end script

post-stop script

end script
